# wrtbwmon
#!/bin/sh /etc/rc.common

START=91

start(){
	wrtbwmon.sh setup /tmp/usage.db
}

stop(){
	wrtbwmon.sh remove
}
