# wrtbwmon.init
#!/bin/sh /etc/rc.common

START=91

start(){
	wrtbwmon setup /mnt/cifs2/usage.db
}

stop(){
	wrtbwmon remove
}

